# puppetagent/files/etc/profile.d/Z99-puppet_path.sh
### FILE MANAGED BY PUPPET ###


## Puppet 4 installs in /opt so add this to the path.
export PATH=$PATH:/opt/puppetlabs/bin
